# .cursor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server --daemonize no
autostart=true
autorestart=true
stderr_logfile=/var/log/redis.err.log
stdout_logfile=/var/log/redis.out.log

[program:fleet-workers]
command=pnpm pm2 start ecosystem.config.js --no-daemon
directory=/workspace
autostart=true
autorestart=true
stderr_logfile=/var/log/fleet.err.log
stdout_logfile=/var/log/fleet.out.log
